# components/cron/controlPanel.service
[Unit]
Description=web server with heatmaps.
After=network.target

[Service]
User=brunocosta
ExecStart=/home/brunocosta/Downloads/git/controlpanel/bin/controlpanel

[Install]
WantedBy=multi-user.target
